# deployments/dpkg/kafka/support/postinst
#!/bin/sh

update-rc.d -f kafka defaults