# deployments/dpkg/kafka/support/prerm
#!/bin/sh
if [ "$1" = "remove" ]; then
update-rc.d -f kafka remove
fi