# local-auth/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:kratos]
command=/usr/bin/kratos serve --config /etc/config/kratos.yml --dev --watch-courier
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/kratos.log
stderr_logfile=/var/log/supervisor/kratos.err
environment=DSN="%(ENV_DSN)s"
priority=10

[program:hydra]
command=/usr/bin/hydra serve all --config /etc/config/hydra.yml --dev
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/hydra.log
stderr_logfile=/var/log/supervisor/hydra.err
environment=DSN="%(ENV_HYDRA_DSN)s"
priority=20

[program:login-ui]
command=node /app/login-ui/server.js
directory=/app/login-ui
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/login-ui.log
stderr_logfile=/var/log/supervisor/login-ui.err
environment=KRATOS_PUBLIC_URL="http://localhost:4433",KRATOS_ADMIN_URL="http://localhost:4434",HYDRA_ADMIN_URL="http://localhost:4445",PORT="3000"
priority=30
